# integrations/upstream-extracted/niurk-smart-workflow-bundle/niurk-smart-workflow-bundle/systemd/niurk-stack-health.service
[Unit]
Description=Niurk stack health validation
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/opt/niurk-directives/ai-endpoints.env
ExecStart=/opt/niurk-directives/healthcheck.sh
