# integrations/upstream-extracted/niurk-smart-workflow-bundle/niurk-smart-workflow-bundle/systemd/niurk-stack-health.timer
[Unit]
Description=Run Niurk stack health validation periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
